# deploy/systemd/void-saga-backup.service
[Unit]
Description=Void Saga SQLite backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/srv/void-saga/repo
ExecStart=/bin/sh /srv/void-saga/repo/deploy/run-backup.sh
