# deploy/systemd/void-saga-backup.service
[Unit]
Description=Void Saga SQLite backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/srv/void-saga/repo
ExecStart=/usr/bin/docker compose --env-file /srv/void-saga/env/.env.production -f docker-compose.prod.yml --profile backup run --rm backup
